# resources/bobdock.service
[Unit]
Description=BobDock
After=graphical-session.target

[Service]
Type=exec
ExecStart=bobdock

[Install]
WantedBy=graphical-session.target
